# templates/default/init.redis.conf.erb
description "redis server"

start on runlevel [23]
stop on shutdown

respawn

exec sudo -u <%= node['redis']['user'] %> /usr/local/bin/redis-server <%= node['redis']['conf_file'] %>
